# initramfs/files/etc/inittab
# LevitateOS BusyBox init configuration
# TEAM_451: Generated by xtask
# TEAM_460: Changed respawn to wait so 'exit' terminates shell

# System initialization
::sysinit:/bin/echo "LevitateOS (BusyBox) starting..."
::sysinit:/bin/mount -t proc proc /proc
::sysinit:/bin/mount -t sysfs sysfs /sys

# Test scripts at root level (kernel has issues with subdirectory files)
# Run: sh /test-core.sh [phase] to test coreutils
# Run: sh /test.sh to test basic ash functionality

# Start interactive shell
::wait:-/bin/ash

# Handle Ctrl+Alt+Del
::ctrlaltdel:/sbin/reboot

# Shutdown hooks
::shutdown:/bin/echo "System shutting down..."
