# initramfs/files/etc/inittab
# LevitateOS BusyBox init configuration
# TEAM_451: Generated by xtask
# TEAM_460: Changed respawn to wait so 'exit' terminates shell
# TEAM_475: Fixed tty setup using setsid/cttyhack for job control

# System initialization
::sysinit:/bin/echo "LevitateOS (BusyBox) starting..."
::sysinit:/bin/mount -t proc proc /proc
::sysinit:/bin/mount -t sysfs sysfs /sys
::sysinit:/bin/mount -t devtmpfs devtmpfs /dev 2>/dev/null || true

# Start shell with proper controlling tty setup
# getty with -n (no login) and -l (login program = /bin/ash) for auto-login
::respawn:/sbin/getty -n -l /bin/ash 0 ttyS0 vt100

# Handle Ctrl+Alt+Del
::ctrlaltdel:/sbin/reboot

# Shutdown hooks
::shutdown:/bin/echo "System shutting down..."
